# deployment/finma-career-coach.service
[Unit]
Description=FINMA Career Coach Streamlit Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/finma-career-coach
Environment="PATH=/var/www/finma-career-coach/venv/bin"
Environment="ANTHROPIC_API_KEY=your_api_key_here"
ExecStart=/var/www/finma-career-coach/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
